# staging_dir/target-mips_r2_uClibc-0.9.33.2/root-ar71xx/etc/init.d/rild
#!/bin/sh /etc/rc.common
# Copyright (C) 2008-2012 OpenWrt.org

START=12

start() 
{	
	service_start /usr/bin/rild
}

stop() 
{
	service_stop /usr/bin/rild
}
