# staging_dir/target-mips_r2_uClibc-0.9.33.2/root-ar71xx/etc/init.d/slicMgr
#!/bin/sh /etc/rc.common
# Copyright (C) 2008-2012 OpenWrt.org

START=71

start() 
{	
	service_start /bin/slicMgr
}

stop() 
{
	service_stop /bin/slicMgr
}
